# files/etc/systemd/system/uos-discovery-client.service.d/10-wait-for-network-setup.conf
[Unit]
After=network_setup.service
Requires=network_setup.service
